# systemd/hop.service
[Unit]
Description=Hop SSH jump server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hop
Group=hop
WorkingDirectory=/var/lib/hop
ExecStart=/usr/local/bin/hop-server serve --config /etc/hop/config.toml
Restart=on-failure
RestartSec=2
Environment=RUST_LOG=info
Environment=HOP_ADMIN_STATIC_DIR=/usr/share/hop/admin-web

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/hop

[Install]
WantedBy=multi-user.target
